# luci-app-halsamrec/root/etc/init.d/halsamrec-luci
#!/bin/sh /etc/rc.common

START=95

start() {
	echo "ALSA device LuCI module initialized"
}

stop() {
	echo "ALSA device LuCI module stopped"
}
