# luci-app-halsamrec/root/etc/init.d/halsamrec-luci
#!/bin/sh /etc/rc.common
# Passive registration — no daemon is managed by this init script

START=95

start() { return 0; }
stop()  { return 0; }
